# backend/crontab
* 0 * * * /usr/local/bin/python /main.py >> /logfile

